# venom_lambda/systemd/venom-mesh-discovery.service
[Unit]

Description=VENOM Λ-Core Mesh Discovery Daemon (UDP Multicast)
Documentation=https://github.com/manuelstellian-dev/VENOM-Omega-Lambda
After=network-online.target
Wants=network-online.target

[Service]

Type=simple
User=%u
WorkingDirectory=%h/venom/lambda/core
ExecStart=%h/.local/bin/python3 mesh_discovery.py
Restart=on-failure
RestartSec=3

# Resource limits
MemoryMax=256M
CPUQuota=25%

# Logging
StandardOutput=append:%h/venom/logs/mesh-discovery.log
StandardError=append:%h/venom/logs/mesh-discovery.error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=VENOM_HOME=%h/venom

# Network
# Allow multicast (required for UDP multicast)
RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=default.target
