# venom_lambda/systemd/venom-mesh-orchestrator.service
[Unit]

Description=VENOM Λ-Core Mesh Orchestrator
Documentation=https://github.com/manuelstellian-dev/VENOM-Omega-Lambda
After=venom-mesh-discovery.service venom-api.service
Requires=venom-mesh-discovery.service
Wants=venom-api.service

[Service]

Type=simple
User=%u
WorkingDirectory=%h/venom/lambda/core
ExecStart=%h/.local/bin/python3 venom-mesh-orchestrator.py
Restart=on-failure
RestartSec=5

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=append:%h/venom/logs/orchestrator.log
StandardError=append:%h/venom/logs/orchestrator.error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=VENOM_HOME=%h/venom
Environment=VENOM_HEALTH_INTERVAL=5

[Install]
WantedBy=default.target
