# venom_lambda/systemd/venom-fractal.service
[Unit]
Description=VENOM Fractal Time Compression Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/VENOM-Omega-Lambda
ExecStart=/usr/bin/python3 -m venom_lambda.core.fractal
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
